# rootfs/usr/lib/opkg/info/dropbear.list
/lib/preinit/99_10_failsafe_dropbear
/usr/bin/dbclient
/etc/dropbear/dropbear_ed25519_host_key
/etc/dropbear/dropbear_rsa_host_key
/etc/init.d/dropbear
/usr/bin/dropbearkey
/etc/config/dropbear
/usr/sbin/dropbear
