# rootfs/usr/lib/opkg/info/firewall4.list
/etc/nftables.d/10-custom-filter-chains.nft
/usr/share/firewall4/templates/zone-mssfix.uc
/lib/upgrade/keep.d/firewall4
/etc/nftables.d/README
/usr/share/firewall4/helpers
/usr/share/firewall4/templates/zone-masq.uc
/etc/init.d/firewall
/usr/share/firewall4/templates/mangle-rule.uc
/sbin/fw3
/sbin/fw4
/etc/config/firewall
/usr/share/nftables.d/README
/usr/share/firewall4/main.uc
/usr/share/ucode/fw4.uc
/usr/share/firewall4/templates/redirect.uc
/usr/share/firewall4/templates/zone-verdict.uc
/usr/share/firewall4/templates/rule.uc
/etc/hotplug.d/iface/20-firewall
/usr/share/firewall4/templates/ruleset.uc
/usr/share/firewall4/templates/zone-drop-invalid.uc
/usr/share/firewall4/templates/zone-notrack.uc
/usr/share/firewall4/templates/zone-jump.uc
/usr/share/firewall4/templates/zone-match.uc
