# rootfs/usr/lib/opkg/info/luci-base.list
/www/luci-static/resources/firewall.js
/etc/config/luci
/www/luci-static/resources/rpc.js
/www/luci-static/resources/cbi/folder.svg
/www/luci-static/resources/icons/alias.png
/usr/share/ucode/luci/controller/admin/uci.uc
/www/luci-static/resources/form.js
/www/luci-static/resources/icons/alias_disabled.png
/www/luci-static/resources/icons/switch_disabled.png
/etc/init.d/ucitrack
/www/luci-static/resources/icons/loading.gif
/usr/share/ucode/luci/controller/admin/index.uc
/usr/share/ucode/luci/template/sysauth.ut
/www/luci-static/resources/icons/signal-75-100.png
/www/luci-static/resources/icons/signal-none.png
/www/luci-static/resources/tools/views.js
/sbin/luci-reload
/www/luci-static/resources/ui.js
/www/luci-static/resources/icons/wifi.png
/usr/share/ucode/luci/template/view.ut
/usr/share/ucode/luci/sys.uc
/www/luci-static/resources/network.js
/www/luci-static/resources/tools/prng.js
/www/luci-static/resources/fs.js
/www/luci-static/resources/icons/port_up.png
/www/luci-static/resources/protocol/dhcp.js
/etc/config/ucitrack
/www/luci-static/resources/promis.min.js
/usr/share/ucode/luci/uhttpd.uc
/usr/share/ucode/luci/template/error500.ut
/www/luci-static/resources/icons/port_down.png
/usr/share/acl.d/luci-base.json
/www/luci-static/resources/icons/tunnel.png
/www/luci-static/resources/xhr.js
/www/luci-static/resources/icons/vlan_disabled.png
/www/luci-static/resources/icons/vlan.png
/lib/upgrade/keep.d/luci-base
/www/luci-static/resources/uci.js
/usr/lib/ucode/luci/core.so
/usr/share/ucode/luci/template/footer.ut
/www/luci-static/resources/cbi.js
/www/luci-static/resources/cbi/file.svg
/www/cgi-bin/luci
/usr/share/ucode/luci/dispatcher.uc
/www/luci-static/resources/icons/signal-50-75.png
/usr/share/ucode/luci/version.uc
/usr/share/luci/menu.d/luci-base.json
/etc/luci-uploads/.placeholder
/usr/share/ucode/luci/template/error404.ut
/usr/share/ucode/luci/template/csrftoken.ut
/usr/share/ucode/luci/template/header.ut
/www/luci-static/resources/tools/widgets.js
/www/luci-static/resources/protocol/none.js
/usr/share/ucode/luci/runtime.uc
/www/luci-static/resources/icons/wifi_disabled.png
/www/luci-static/resources/icons/signal-0-25.png
/usr/share/rpcd/acl.d/luci-base.json
/www/luci-static/resources/protocol/static.js
/www/luci-static/resources/validation.js
/www/luci-static/resources/icons/tunnel_disabled.png
/www/luci-static/resources/icons/signal-0.png
/www/luci-static/resources/icons/switch.png
/usr/share/rpcd/ucode/luci
/www/luci-static/resources/icons/ethernet_disabled.png
/www/luci-static/resources/icons/bridge_disabled.png
/www/luci-static/resources/cbi/link.svg
/www/luci-static/resources/luci.js
/usr/share/ucode/luci/http.uc
/usr/share/ucode/luci/zoneinfo.uc
/www/luci-static/resources/icons/ethernet.png
/www/index.html
/www/luci-static/resources/icons/signal-25-50.png
/www/luci-static/resources/icons/bridge.png
