# rootfs/usr/lib/opkg/info/netifd.list
/lib/netifd/proto/dhcp.sh
/etc/init.d/network
/lib/netifd/netifd-proto.sh
/etc/hotplug.d/iface/00-netstate
/sbin/ifdown
/lib/netifd/dhcp.script
/usr/share/udhcpc/default.script
/lib/network/config.sh
/lib/upgrade/keep.d/netifd
/sbin/netifd
/sbin/devstatus
/etc/init.d/packet_steering
/lib/netifd/utils.sh
/lib/netifd/netifd-wireless.sh
/etc/udhcpc.user
/usr/libexec/network/packet-steering.sh
/sbin/ifstatus
/sbin/ifup
/etc/uci-defaults/14_migrate-dhcp-release
