# rootfs/usr/lib/opkg/info/odhcpd-ipv6only.list
/usr/sbin/odhcpd
/etc/init.d/odhcpd
/usr/sbin/odhcpd-update
/etc/uci-defaults/15_odhcpd
