# rootfs/usr/lib/opkg/info/ttyd.list
/usr/bin/ttyd
/etc/init.d/ttyd
/etc/config/ttyd
